# .config/systemd/user/wlsunset.service
[Unit]
Description=Day/Night Gamma Adjustment
After=multi-user.target

[Service]
ExecStart=/usr/bin/wlsunset -l -37.8 -L 145
