# dotfiles/.config/systemd/user/mailsync.service
[Unit]
Description=Local IMAP Sync

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mailsync
