# dotfiles/.config/systemd/user/mailsync.service
[Unit]
Description=Mailsync

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mailsync
